# com.roo.trading-report-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.roo.trading-report-push</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/andrewbuis/Claude/Projects/trading-info-page/push-report.sh</string>
  </array>

  <!-- Fires twice daily:
       12:15 AM — 15 min after Claude writes the overnight/daily report at midnight
       6:15 PM  — 15 min after Claude writes the EOD assessment at 6 PM -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>0</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>15</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/andrewbuis/Claude/Projects/trading-info-page/push.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/andrewbuis/Claude/Projects/trading-info-page/push.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
